# .config/service/pipewire/run
#!/usr/bin/env sh

exec /usr/bin/pipewire 2>&1
