# .config/service/emacs/run
#!/usr/bin/env sh

exec /usr/bin/emacs --fg-daemon 2>&1
